# user/.crontab
MAILTO=""
@daily /usr/local/bin/brewgrade
@daily /usr/local/bin/youtube-dl -U
@hourly /usr/local/bin/wallpaper
